# meta-custom/meta-bachelor/recipes-userspace/serverapp/files/serverapp.service
[Unit]
Description=serverapp
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/serverapp

[Install]
WantedBy=multi-user.target
